# systemd/trident.service
[Unit]
Description=Trident Agent
After=network.target network-online.target systemd-udev-settle.service

[Service]
ExecStart=trident commit --verbosity DEBUG
Type=oneshot

[Install]
WantedBy=multi-user.target