# packaging/systemd/tridentd.service
[Unit]
Description=Trident Daemon

[Service]
ExecStart=trident daemon
